# infra/backups/systemd/backup-cleanup.timer
[Unit]
Description=Asya backup retention cleanup timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=backup-cleanup.service

[Install]
WantedBy=timers.target
