# infra/backups/systemd/pg-backup-hourly.timer
[Unit]
Description=Asya PostgreSQL incremental backup (hourly)

[Timer]
OnCalendar=hourly
Persistent=true
Unit=pg-backup@incremental.service

[Install]
WantedBy=timers.target
