# infra/backups/systemd/s3-offsite-sync.timer
[Unit]
Description=Asya object storage off-site sync timer

[Timer]
OnCalendar=hourly
Persistent=true
Unit=s3-offsite-sync.service

[Install]
WantedBy=timers.target
